# conf/sova.conf
[program:sova]
command=/sova/target/release/sova
