# sample_app/sample_app.conf
[program:sample_app]
command=/sample_app/main.py
